# packaging/ember-llm.service
[Unit]
Description=EmberOS LLM Servers (BitNet + Qwen2.5-VL)
Documentation=https://docs.emberos.org
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ember-llm-manager

Restart=on-failure
RestartSec=10

# Environment
Environment=LLAMA_LOG_LEVEL=INFO

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/usr/local/share/ember/models
PrivateTmp=true

[Install]
WantedBy=default.target

