# packaging/emberd.service
[Unit]
Description=EmberOS Agent Daemon
Documentation=https://docs.emberos.org
After=dbus.service ember-llm.service
Wants=ember-llm.service

[Service]
Type=simple
ExecStart=%h/.local/share/ember/venv/bin/emberd
Restart=on-failure
RestartSec=5

# Environment
Environment=EMBER_LOG_LEVEL=INFO
Environment=CHROMA_CACHE_DIR=%h/.local/share/ember/vectors/.cache

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/ember %h/.config/ember %h/.cache/ember %h/Downloads %h/Documents %h/Desktop %h/Pictures %h/Videos %h/Music %h
PrivateTmp=false

[Install]
WantedBy=default.target
